# profile.d/bcompare.guard.sh
bcompare() (
    QT_GRAPHICSSYSTEM=native command bcompare "$@"
)
f.complete bcompare

bcompare.session() {
    :
}
f.complete bcompare.session
bcompare.session


