# profile.d/cargo.source.sh
# see also cargo.guard.sh, which does further configuration

source.if ~/.cargo/env
