# profile.d/cargo.source.sh
# assume path.sh sourced
path.add ~/.cargo/bin



