# profile.d/copyq.guard.sh
copyq() (
     QT_QPA_PLATFORM=xcb command copyq "$@"
)
f.complete copyq

copyq.session() (
    :
)
f.complete copyq.session
copyq.session

