# profile.d/go.guard.sh
# assume path.sh sourced
# path.add $(go env GOPATH)/bin
export GOBIN=${HOME}/.local/bin


