# profile.d/never.guard.sh
>&2 echo "${BASH_SOURCE} should never be guarded."
u.have ${BASH_SOURCE%%.guard.sh}
