# profile.d/packer.source.sh
packer.session() {
    complete -C packer packer
}
f.x packer.session
packer.session

loaded "${BASH_SOURCE}"   

