# profile.d/torbrowser-launcher.source.sh
tor.downloads() { cd ~/.local/share/torbrowser/tbb/x86_64/tor-browser/Browser/Downloads/; }
f.x tor.downloads

