# profile.d/.config/systemd/user/mnt-gdrive-mcarifio.automount
[Unit]
Description=mount rclone %p for %u
After=network-online.target

[Service]
User=%u
Group=%g

ExecStartPre=mkdir -p %h/mnt/gdrive/carif.io || true

ExecStart=/usr/bin/rclone mount %p:/ %h/mnt/gdrive/carif.io --vfs-cache-mode full --drive-export-formats docx --drive-export-formats xlsx --drive-import-formats docx --drive-import-formats xlsx    
